# DataCollector/crontab
0 6 * * * /usr/local/bin/python /app/apiCollector.py >> /var/log/cron.log 2>&1
